# host/box-firewall.service
[Unit]
Description=Re-apply box host firewall rules
# After docker so DOCKER-USER exists; still fine if docker is absent.
After=network.target incus.service docker.service
Wants=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/box-firewall

[Install]
WantedBy=multi-user.target
